# password-factory.service.in
[Unit]
Description=Password Factory

[Service]
User=@USER@
Group=@GROUP@
ExecStart=@bindir@/password-factory \
	--ssl-certificate-file=@ssl_cert_file@ \
	--ssl-private-key-file=@ssl_key_file@ \
	--template-directory=@templdir@ \
	--document-root=@basedir@

[Install]
WantedBy=multi-user.target
